# crontab
* 20 * * * python /home/app/power_outage_notification.py >> /logfile